# docker/crontab
* * * * * root /usr/bin/feminicides-update /usr/share/nginx/html/data current
